# build/build/cfg/etc/cron.hourly/slrnpull
#!/bin/sh
umask 022
slrnpull -d /var/spool/slrnpull -h news.eternal-september.org
