# build/build/root/etc/sv/syncthing/run
#!/bin/sh

export USERNAME=shum
export HOME="/home/$USERNAME"
export SYNCTHING="/usr/bin/syncthing"

exec 2>&1
exec chpst -u "$USERNAME" "$SYNCTHING" -no-browser -no-restart -logflags 0

