# kobo/KoboStuff/src/etc/udev/rules.d/99-stuff.rules
# $Id: 99-stuff.rules 11379 2015-01-10 23:58:00Z NiLuJe $
# Runs early at boot... (onboard *might* be mounted at that point)
KERNEL=="loop0", RUN+="/usr/local/stuff/bin/stuff.sh"
