# build/build/cfg/etc/sv/sonarr/run
#!/bin/sh
export SONARR=/srv/http/sonarr
export USERNAME=sonarr
exec 2>&1
exec chpst -u "$USERNAME:$USERNAME" mono /usr/lib/sonarr/NzbDrone/NzbDrone.exe -s localhost:8989 -no-browser -data=$SONARR >/dev/null

