# build/build/cfg/etc/sv/dhcpcd-eno1/run
#!/bin/sh
[ -r conf ] && . ./conf
exec dhcpcd -B eno1 $OPTS 1>&2
