# build/build/cfg/etc/sv/pulseaudio-browser/run
#!/bin/sh
export USERNAME=shum
exec 2>&1
# browser doesn't autostart pulseaudio
exec chpst -u "$USERNAME" /usr/bin/pulseaudio --start --log-target=syslog
