# build/build/cfg/etc/sv/php-fpm/run
#!/bin/sh
install -d -m 0755 -o root -g root /var/run/php-fpm
exec 2>&1
exec php-fpm --nodaemonize
