# build/build/cfg/etc/sv/nzbget/run
#!/bin/sh
export USERNAME=$USER
export NZBGET="/usr/bin/nzbget"
exec 2>&1
exec chpst -u "$USERNAME" "$NZBGET" --daemon >/dev/null

