# build/build/cfg/etc/sv/nzbget/run
#!/bin/sh
export USERNAME=shum
exec 2>&1
# runit does not detect bg (--daemon) process, use --server --option to spawn single fg process
exec chpst -u "$USERNAME" /usr/bin/nzbget --server --option OutputMode=loggable -c /home/$USERNAME/.nzbget
